# share/etc/services/system/rdm.socket
[Unit]
Description=RTags daemon socket

[Socket]
ListenStream=%t/rdm.socket

[Install]
WantedBy=default.target
